# deploy/systemd/hpcdeploy-frontend.service
[Unit]
Description=HPCDeploy Vue Frontend (Vite Dev Server)
After=network.target hpcdeploy-backend.service

[Service]
Type=simple
User=REPLACED_BY_INSTALL_SCRIPT
Group=REPLACED_BY_INSTALL_SCRIPT
WorkingDirectory=REPLACED_BY_INSTALL_SCRIPT/frontend
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0 --port 5173
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
